# home/Library/LaunchAgents/com.user.sing-box-rulesets.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.sing-box-rulesets</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jacob/.config/sing-box/update-rulesets.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>17</integer>
        <key>Weekday</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/opt/homebrew/var/log/sing-box/rulesets-update.log</string>
    <key>StandardErrorPath</key>
    <string>/opt/homebrew/var/log/sing-box/rulesets-update.log</string>
</dict>
</plist>
